# CLONE/Twilio_0/Install/supervisor.conf
# This file is in /etc/supervisord.d
[program:agegroupltd]
command=gunicorn -c /root/SMSTask/CLONE/Twilio_0/gunicorn.conf.py Twilio.wsgi:application
directory=/root/SMSTask/CLONE/Twilio_0
user=root
autorestart=true
startretires=3
stderr_logfile=/root/SMSTask/CLONE/Twilio_0/log/supervisor.err.log
stdout_logfile=/root/SMSTask/CLONE/Twilio_0/log/supervisor.out.log